# root/etc/services.d/myjd_api/run
#!/usr/bin/with-contenv bash

exec myjd_api --docker $PARAMETER
fi
